# gentooplayer/script/ramdisk.init
#!/sbin/openrc-run

user="root:root"
logfile="/var/log/ramdisk.log"

start_stop_daemon_args="--user $user"

command="/etc/default/ramdisk.sh"
command_args="
	-f $logfile
"

command_background=no
#pidfile=/run/roon.pid

depend() {
    need localmount
    after bootmisc
    before xdm
}

start_pre() {
    checkpath --file --owner $user --mode 0644 $logfile
}
